# aprphpfpm/etc/services.d/php-fpm/finish
#!/bin/sh

echo '... ... Service php-fpm[php-fpm7] Finished ... ...'
