# aprphpfpm/etc/services.d/php-fpm/finish
#!/bin/sh

pkill -SIGTERM php-fpm7
echo '... ... Service php-fpm[php-fpm7] Finished ... ...'
