# aprphpfpm/etc/services.d/php-fpm/run
#!/bin/sh

sux www-data php-fpm7 --nodaemonize --fpm-config /etc/php7/php-fpm.conf

echo '... ... Service php-fpm[php-fpm7] Started Successfully ... ...'
