# s6version/apsapache/etc/services.d/apache2-access-log/run
#!/usr/bin/execlineb -P
logutil-service -f /var/run/s6/apache2-access-log-fifo /var/log/apache2/access-logs
