# aprwpbase/etc/services.d/php-error-log/run
#!/usr/bin/execlineb -P
logutil-service -f /var/run/s6/php-error-log-fifo /var/log/php/error-logs
