# aprnginx/etc/services.d/nginx/finish
#!/bin/sh

echo '... ... Service nginx Finished ... ...'
