# apsphpfpm/etc/services.d/php-access-log/run
#!/usr/bin/execlineb -P
logutil-service -f /var/run/s6/php-access-log-fifo /var/log/php/access-logs
