# s6version/apsapache/etc/services.d/apache2-error-log/run
#!/usr/bin/execlineb -P
logutil-service -f /var/run/s6/apache2-error-log-fifo /var/log/apache2/error-logs
